# dist/cluegetter.service
[Unit]
Description=Access and Auditing Milter
After=nss-lookup.target

[Service]
Restart=always
RestartSec=30
EnvironmentFile=/etc/default/cluegetter
ExecStartPre=-/bin/mkdir -p /var/run/cluegetter/
ExecStartPre=-/bin/chown -R cluegetter /var/run/cluegetter/
ExecStart=/usr/bin/cluegetter --config /etc/cluegetter/cluegetter.conf --loglevel $LOGLEVEL daemon
# User=cluegetter

[Install]
WantedBy=multi-user.target
